# dotfiles/qtile/dot-config/autostart/dunst.desktop
[Desktop Entry]
Name=Dunst
Comment=Dunst is a lightweight notification daemon
Exec=dunst
Terminal=false
Type=Application
