# dotfiles/qtile/dot-config/autostart/kanshi.desktop
[Desktop Entry]
Name=Kanshi
Comment=Kanshi is a Wayland display configuration daemon
Exec=kanshi
Terminal=false
Type=Application
